# rover/controls-systems/mobility/GNSS/PiControl/reach.service
[Unit]
Description=Reach Sub System

[Service]
ExecStart=/etc/systemd/system/reach_USB_OTG.sh
User=pi
Group=pi

[Install]
WantedBy=multi-user.target
